# ops/systemd/skibot-cleanup-prod.timer
[Unit]
Description=Skibot PROD cleanup — every hour

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
Unit=skibot-cleanup-prod.service

[Install]
WantedBy=timers.target
